# docker/supervisord.conf
; /etc/supervisor/conf.d/supervisord.conf

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; 后端服务
[program:backend]
; 使用单worker模式避免多进程问题，增加超时设置
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --timeout-keep-alive 300 --limit-max-requests 1000
directory=/unified-app
environment=PYTHONPATH="/unified-app"
autostart=true
autorestart=true
priority=10
startretries=10
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend-error.log
startsecs=3
; 增加停止等待时间，确保请求处理完成
stopwaitsecs=30
stopsignal=TERM

; Nginx服务
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
startsecs=2

; 进程组管理
[group:app]
programs=backend,nginx
